# supervisor.ini
[inet_http_server]
port = :9001

[supervisord]
procFile = /tmp/proc.json
loglevel= trace
logfile=/dev/stdout

[program:mytest]
directory = /root/
command = nc -l 8090

